# extras/satania.timer
[Unit]
Description=Satania Timer

[Timer]
OnCalendar=*:0/20

[Install]
WantedBy=timers.target
